# Goblin/Stager10/fluend.service
[Unit]
Description=fluend

[Service]
ExecStart=/gbln/check-life.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target 
